# debos/overlays/007-splash-screen/home/neon/.config/autostart/set-wallpaper.desktop
[Desktop Entry]
Type=Application
Name=Set Default Wallpaper
Exec=/usr/local/bin/set-default-wallpaper.sh
Hidden=false
X-KDE-autostart-condition=rcfile
