# debos/overlays/007-splash-screen/home/neon/.config/systemd/user/wallpaper.service
[Unit]
Description=Set Neon splash screen as user wallpaper
DefaultDependencies=no
After=plasma-workspace.target
Requires=plasma-workspace.target

[Service]
Type=oneshot
ExecStart=/usr/bin/plasma-apply-wallpaperimage /opt/neon/splashscreen.png
RemainAfterExit=yes
StandardOutput=journal+console
StandardError=journal+console

[Install]
WantedBy=plasma-workspace.target